# config/supervisor/kubetty.conf
; KubeTTY Supervisor Configuration
; Manages kubetty-project and GUI stack (Xvfb, x11vnc, XFCE)
;
; GUI components only start when GUI_ENABLED=true environment variable is set.

[supervisord]
nodaemon=true
; Note: No 'user=' directive - runs as container user (mmattox)
; This allows supervisor to work in non-root containers
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=10MB
logfile_backups=3

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; =============================================================================
; KubeTTY Project Server (always runs)
; =============================================================================
[program:kubetty-project]
command=/usr/local/bin/kubetty-project
; user directive removed - runs as container user (mmattox)
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/kubetty-project.log
stderr_logfile=/var/log/supervisor/kubetty-project-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; =============================================================================
; GUI Stack (only starts when GUI_ENABLED=true)
; =============================================================================

; GUI Initializer - starts GUI components in correct order
; This program runs start-gui.sh which handles XDG setup and sequential startup
[program:gui-init]
command=/usr/local/bin/start-gui
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=15
stdout_logfile=/var/log/supervisor/gui-init.log
stderr_logfile=/var/log/supervisor/gui-init-error.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB

; Xvfb - X Virtual Framebuffer (headless X server)
; Creates display :99 with configured resolution
; GUI_RESOLUTION format: WIDTHxHEIGHTxDEPTH (e.g., 1920x1080x24)
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_GUI_RESOLUTION)s -ac +extension GLX +render -noreset
autostart=false
autorestart=true
priority=20
environment=HOME="/home/mmattox"
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; x11vnc - VNC server exposing Xvfb display
; GUI_VNC_PORT: VNC port (5900-5999), defaults to 5900
; Note: -ncache is NOT used because noVNC doesn't support it (expands framebuffer)
; Anti-blanking options:
;   -sb 0: Disable screen blank throttling (default 60s causes "freeze" appearance)
;   -nodpms: Prevent DPMS power management when clients connected
;   -nofbpm: Prevent frame buffer power management
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -rfbport %(ENV_GUI_VNC_PORT)s -xkb -sb 0 -nodpms -nofbpm
autostart=false
autorestart=true
priority=30
startsecs=3
depends_on=xvfb
environment=DISPLAY=":99",HOME="/home/mmattox"
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; D-Bus session daemon (required for XFCE)
[program:dbus]
command=/usr/bin/dbus-daemon --session --nofork --print-address
autostart=false
autorestart=true
priority=25
environment=DISPLAY=":99",HOME="/home/mmattox"
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus-error.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB

; XFCE Desktop Session
[program:xfce]
command=/usr/bin/startxfce4
autostart=false
autorestart=true
priority=40
startsecs=5
depends_on=xvfb,dbus
environment=DISPLAY=":99",HOME="/home/mmattox",XDG_RUNTIME_DIR="/run/user/1000"
stdout_logfile=/var/log/supervisor/xfce.log
stderr_logfile=/var/log/supervisor/xfce-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; =============================================================================
; Program Groups
; =============================================================================
[group:gui]
programs=xvfb,dbus,x11vnc,xfce
priority=100
